# config/busybox-1.22.1/files
# vim: syntax=sh

# Get the busybox-built binary files.
bin/busybox
'symlink' bin/*      busybox
'symlink' sbin/*     ../bin/busybox
'symlink' usr/bin/*  ../../bin/busybox
'symlink' usr/sbin/* ../../bin/busybox
init
linuxrc

# Get the rootfs overlay files.
etc/fstab
etc/group
etc/hostapd.conf
etc/hostname
etc/hosts
etc/httpd.conf
etc/inittab
etc/inputrc
etc/issue
etc/motd
etc/mtab
etc/passwd
etc/profile
etc/protocols
etc/rc.d/init.d/ftpd
etc/rc.d/init.d/hostapd
etc/rc.d/init.d/httpd
etc/rc.d/init.d/network
etc/rc.d/init.d/random
etc/rc.d/init.d/syslog
etc/rc.d/init.d/telnetd
etc/rc.d/init.d/udhcpd
etc/rc.d/rc.initscript
etc/rc.d/rc.shutdown/10.udhcpd
etc/rc.d/rc.shutdown/20.telnetd
etc/rc.d/rc.shutdown/30.htpd
etc/rc.d/rc.shutdown/40.ftpd
etc/rc.d/rc.shutdown/60.hostapd
etc/rc.d/rc.shutdown/70.network
etc/rc.d/rc.shutdown/80.syslog
etc/rc.d/rc.shutdown/90.random
etc/rc.d/rc.startup/10.random
etc/rc.d/rc.startup/20.syslog
etc/rc.d/rc.startup/30.network
etc/rc.d/rc.startup/40.hostapd
etc/rc.d/rc.startup/60.ftpd
etc/rc.d/rc.startup/70.httpd
etc/rc.d/rc.startup/80.telnetd
etc/rc.d/rc.startup/90.udhcpd
etc/rc.d/rc.sysdone
etc/rc.d/rc.sysinit
etc/securetty
etc/services
etc/shadow
etc/shells
etc/udhcpd.conf
etc/wpa_supplicant.conf
home/ftpd
home/httpd
home/httpd/index.html
home/tftpd
home/user
usr/share/udhcpc/default.script
var/log/lastlog
