# config/busybox-1.22.1/files
# vim: syntax=sh
bin/ash
bin/awk
bin/base64
bin/busybox
bin/cat
bin/catv
bin/chattr
bin/chgrp
bin/chmod
bin/chown
bin/cp
bin/cpio
bin/date
bin/dd
bin/df
bin/dmesg
bin/dnsdomainname
bin/dumpkmap
bin/echo
bin/egrep
bin/false
bin/fgrep
bin/getopt
bin/grep
bin/gunzip
bin/gzip
bin/hostname
bin/kill
bin/ln
bin/login
bin/ls
bin/lsattr
bin/mkdir
bin/mknod
bin/mktemp
bin/more
bin/mount
bin/mountpoint
bin/mv
bin/netstat
bin/nice
bin/pidof
bin/ping
bin/pipe_progress
bin/printenv
bin/ps
bin/pwd
bin/rm
bin/rmdir
bin/run-parts
bin/sed
bin/setserial
bin/sh
bin/sleep
bin/stat
bin/stty
bin/su
bin/sync
bin/tar
bin/touch
bin/true
bin/umount
bin/uname
bin/usleep
bin/vi
bin/watch
bin/zcat
etc/fstab
etc/group
etc/hostapd.conf
etc/hostname
etc/hosts
etc/httpd.conf
etc/inittab
etc/inputrc
etc/issue
etc/motd
etc/mtab
etc/passwd
etc/profile
etc/protocols
etc/rc.d/init.d/ftpd
etc/rc.d/init.d/hostapd
etc/rc.d/init.d/httpd
etc/rc.d/init.d/network
etc/rc.d/init.d/random
etc/rc.d/init.d/syslog
etc/rc.d/init.d/telnetd
etc/rc.d/init.d/udhcpd
etc/rc.d/rc.initscript
etc/rc.d/rc.shutdown/10.udhcpd
etc/rc.d/rc.shutdown/20.telnetd
etc/rc.d/rc.shutdown/30.htpd
etc/rc.d/rc.shutdown/40.ftpd
etc/rc.d/rc.shutdown/60.hostapd
etc/rc.d/rc.shutdown/70.network
etc/rc.d/rc.shutdown/80.syslog
etc/rc.d/rc.shutdown/90.random
etc/rc.d/rc.startup/10.random
etc/rc.d/rc.startup/20.syslog
etc/rc.d/rc.startup/30.network
etc/rc.d/rc.startup/40.hostapd
etc/rc.d/rc.startup/60.ftpd
etc/rc.d/rc.startup/70.httpd
etc/rc.d/rc.startup/80.telnetd
etc/rc.d/rc.startup/90.udhcpd
etc/rc.d/rc.sysdone
etc/rc.d/rc.sysinit
etc/securetty
etc/services
etc/shadow
etc/shells
etc/udhcpd.conf
etc/wpa_supplicant.conf
home/ftpd
home/httpd
home/httpd/index.html
home/tftpd
home/user
init
linuxrc
sbin/arp
sbin/blockdev
sbin/devmem
sbin/fsck
sbin/getty
sbin/halt
sbin/ifconfig
sbin/ifdown
sbin/ifup
sbin/init
sbin/inotifyd
sbin/insmod
sbin/ip
sbin/ipaddr
sbin/iplink
sbin/iproute
sbin/iprule
sbin/iptunnel
sbin/klogd
sbin/loadkmap
sbin/losetup
sbin/lsmod
sbin/makedevs
sbin/mdev
sbin/mkswap
sbin/modinfo
sbin/modprobe
sbin/nameif
sbin/pivot_root
sbin/poweroff
sbin/reboot
sbin/rmmod
sbin/route
sbin/runlevel
sbin/setconsole
sbin/start-stop-daemon
sbin/sulogin
sbin/swapoff
sbin/swapon
sbin/switch_root
sbin/sysctl
sbin/syslogd
sbin/udhcpc
sbin/vconfig
usr/bin/[
usr/bin/[[
usr/bin/basename
usr/bin/bunzip2
usr/bin/bzcat
usr/bin/bzip2
usr/bin/chrt
usr/bin/chvt
usr/bin/cksum
usr/bin/clear
usr/bin/cmp
usr/bin/crontab
usr/bin/cut
usr/bin/dc
usr/bin/deallocvt
usr/bin/diff
usr/bin/dirname
usr/bin/dos2unix
usr/bin/du
usr/bin/dumpleases
usr/bin/env
usr/bin/expr
usr/bin/find
usr/bin/fold
usr/bin/free
usr/bin/ftpget
usr/bin/ftpput
usr/bin/fuser
usr/bin/head
usr/bin/hexdump
usr/bin/hostid
usr/bin/id
usr/bin/install
usr/bin/ipcrm
usr/bin/ipcs
usr/bin/killall
usr/bin/last
usr/bin/less
usr/bin/logger
usr/bin/logname
usr/bin/lsof
usr/bin/md5sum
usr/bin/mesg
usr/bin/microcom
usr/bin/mkfifo
usr/bin/nc
usr/bin/nohup
usr/bin/nslookup
usr/bin/od
usr/bin/openvt
usr/bin/passwd
usr/bin/patch
usr/bin/printf
usr/bin/readlink
usr/bin/realpath
usr/bin/renice
usr/bin/reset
usr/bin/resize
usr/bin/seq
usr/bin/setkeycodes
usr/bin/setsid
usr/bin/sha1sum
usr/bin/sha256sum
usr/bin/sha3sum
usr/bin/sha512sum
usr/bin/sort
usr/bin/strings
usr/bin/tail
usr/bin/tcpsvd
usr/bin/tee
usr/bin/telnet
usr/bin/test
usr/bin/time
usr/bin/top
usr/bin/tr
usr/bin/traceroute
usr/bin/tty
usr/bin/udpsvd
usr/bin/uniq
usr/bin/unix2dos
usr/bin/unzip
usr/bin/uptime
usr/bin/uudecode
usr/bin/uuencode
usr/bin/vlock
usr/bin/wc
usr/bin/wget
usr/bin/which
usr/bin/who
usr/bin/whoami
usr/bin/whois
usr/bin/xargs
usr/bin/yes
usr/sbin/addgroup
usr/sbin/adduser
usr/sbin/arping
usr/sbin/brctl
usr/sbin/chroot
usr/sbin/crond
usr/sbin/delgroup
usr/sbin/deluser
usr/sbin/dhcprelay
usr/sbin/dnsd
usr/sbin/ether-wake
usr/sbin/ftpd
usr/sbin/httpd
usr/sbin/inetd
usr/sbin/killall5
usr/sbin/ntpd
usr/sbin/rdate
usr/sbin/readprofile
usr/sbin/setlogcons
usr/sbin/telnetd
usr/sbin/tftpd
usr/sbin/udhcpd
usr/share/udhcpc/default.script
var/log/lastlog
var/spool
